# _cron/apa_ping_google.cron
@daily /srv/sites/apa/envs/apa/bin/python /srv/sites/apa/proj/apa/manage.py apa_ping_google
